# sim_deploy_server.crontab
HOME=/home/sim_deploy_server/script
SHELL=/home/sim_deploy_server/venv/bin/python
MAILTO=""
*/5 * * * * ./shop_generator.py
