# sim_deploy_server.crontab
# HOME=/home/sim_deploy_server/script
# SHELL=/home/sim_deploy_server/venv/bin/python
MAILTO=""
*/1 * * * * /home/sim_deploy_server/venv/bin/python /home/sim_deploy_server/script/shop_generator.py
