# sim_deploy_server.crontab
MAILTO=""
*/5 * * * * /home/sim_deploy_server/venv/bin/python /home/sim_deploy_server/script/shop_generator.py;
sleep 60;
/home/sim_deploy_server/venv/bin/python /home/sim_deploy_server/script/upload_to_database.py
