# sim_deploy_server.crontab
SHELL=/home/sim_deploy_server/venv/bin/python
MAILTO=""
*/5 * * * * /home/sim_deploy_server/script/shop_generator.py
