# sim_deploy_server.crontab
HOME=/home/sim_deploy_server/script
SHELL=/home/sim_deploy_server/venv/bin/python
MAILTO=""
*/5 * * * * date; ./shop_generator.py
