# sim_deploy_server.crontab
MAILTO=""
script1='/home/sim_deploy_server/venv/bin/python /home/sim_deploy_server/script/shop_generator.py'
script2='/home/sim_deploy_server/venv/bin/python /home/sim_deploy_server/script/upload_to_database.py'
*/2 * * * * script1 && script2
